# setup/kismet.service
[Unit]
Description=Kismet Server
After=gpsd.service

[Service]
ExecStartPre=/etc/kismet/forcemon.sh wlan1
ExecStart=/usr/bin/kismet_server -s -f /etc/kismet/kismet.conf

Restart=on-abnormal

[Install]
WantedBy=multi-user.target
