# setup/ptuweb.service
[Unit]
Description=PenTestUnit Web Server
After=mysql.service

[Service]
WorkingDirectory=/root/server/src/
ExecStart=/usr/bin/python2 /root/server/src/run.py --host 0.0.0.0 --port 443 --ssl
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target
